# deploy/systemd/aether-webd.service
[Unit]
Description=Aether WebUI Backend API Service
Documentation=https://github.com/bengrewell/aether-webui
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aether-webd
Group=aether-webd

# State directory for persistent data (database)
StateDirectory=aether-webd

# Environment file for additional configuration (must come before ExecStart)
EnvironmentFile=-/etc/aether-webd/env

# Set HOME to the state directory so git does not attempt to read
# /home/aether-webd/.config/git/* (blocked by ProtectHome=true).
Environment=HOME=/var/lib/aether-webd

# Binary location — extra flags come from AETHER_WEBD_OPTS in the env file
ExecStart=/usr/local/bin/aether-webd $AETHER_WEBD_OPTS

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aether-webd

[Install]
WantedBy=multi-user.target
